# .config/systemd/user/hyprlock.service
[Unit]
Description=Hyprland lock screen
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprlock

[Install]
WantedBy=sleep.target
